# config/jarvis.service
[Unit]
Description=Jarvis PTT Daemon
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/Work/jarvis
ExecStart=%h/Work/jarvis/.venv/bin/python -m jarvis.daemon
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
